# configs/leadhub.service
[Unit]
Description=LeadHub Professional API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/leadhub_api
ExecStart=/home/ubuntu/leadhub_api/venv/bin/python /home/ubuntu/leadhub_api/api_integration.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=leadhub_api

[Install]
WantedBy=multi-user.target

